# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── 1. Virtual framebuffer ────────────────────────────────────────────────────
[program:xvfb]
command=Xvfb :1 -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s -ac +extension GLX +render -noreset
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

# ── 2. DBus session (needed by some apps) ────────────────────────────────────
[program:dbus]
command=dbus-daemon --session --address=unix:path=/tmp/dbus-session.sock --nofork --print-address
environment=DISPLAY=":1"
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus.log

# ── 3. Openbox window manager ─────────────────────────────────────────────────
[program:openbox]
command=/bin/bash -c "sleep 1 && openbox-session"
environment=DISPLAY=":1",HOME="/root"
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/openbox.log
stderr_logfile=/var/log/supervisor/openbox.log

# ── 4. x11vnc ─────────────────────────────────────────────────────────────────
[program:x11vnc]
command=/bin/bash -c "sleep 2 && x11vnc -display :1 -rfbauth /root/.vnc/passwd -rfbport 5901 -forever -shared -noxdamage -noxfixes -repeat"
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log

# ── 5. noVNC web bridge on port 3000 ─────────────────────────────────────────
[program:novnc]
command=/bin/bash -c "sleep 3 && websockify --web /usr/share/novnc 3000 localhost:5901"
autorestart=true
priority=50
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.log
